# configure/worker.ini
[program:worker]
command=/home/lepython/work/venv_py3/bin/python /home/lepython/work/venv_py3/bi_system/manage.py celery worker -Q celery -l info -n localhost.celery
process_name=%(program_name)s
numprocs=1
directory=/home/lepython/work/venv_py3/bi_system
umask=022
priority=999
autostart=true
autorestart=unexpected
startsecs=10
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
stopasgroup=false
killasgroup=false
user=lepython
redirect_stderr=false
stdout_logfile=/home/lepython/work/venv_py3/bi_system/run/worker_out.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
stderr_logfile=/home/lepython/work/venv_py3/bi_system/run/worker_err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=1MB
stderr_events_enabled=false
environment=A="1",B="2"
serverurl=AUTO
